# roles/display/templates/monstandby.service.j2
[Unit]
Description=Check monitor alarm standby status
[Service]
ExecStart=/usr/local/bin/monstandby.sh