# roles/display/templates/monstandby.timer.j2
[Unit]
Description=Check monitor alarm standby status
[Timer]
OnBootSec=30
OnUnitActiveSec=8
AccuracySec=1ms
[Install]
WantedBy=timers.target